# devices/MSI-MS-B09611/user-muz/.config/autostart/irexec.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/irexec -d
