# devices/MSI-MS-B09611/user-muz/.config/autostart/setmode.desktop
[Desktop Entry]
Type=Application
Exec=xrandr --output HDMI2  --mode 1280x720
