# devices/MSI-MS-B09611/user-muz/.config/autostart/x11vnc.desktop
[Desktop Entry]
Type=Application
Exec=/home/muz/bin/vnc
Terminal=false
Hidden=false
