# devices/MSI-MS-B09611/user-muz/.config/autostart/xkb-rdp.desktop
[Desktop Entry]
Type=Application
Exec=/home/muz/bin/xkb-rdp.sh
